# templates/etc/systemd/system/caddy.service.erb
[Unit]
Description=Caddy HTTP/2 web server
Documentation=https://caddyserver.com/docs
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]

User=<%= @caddy_user %>
Group=<%= @caddy_group %>

Environment=HOME=<%= @caddy_home %>
Environment=CADDYPATH=<%= @caddy_ssl_dir %>

ExecStart=<%= @install_path %>/caddy -log <%= @caddy_log_dir %>/caddy.log -agree=true -conf=/etc/caddy/Caddyfile -root=/var/tmp -http-port=<%= @caddy_http_port %> -https-port=<%= @caddy_https_port %>
ExecReload=/bin/kill -USR1 $MAINPID

LimitNOFILE=1048576
LimitNPROC=<%= @caddy_limit_processes %>

PrivateTmp=true
PrivateDevices=<%= @caddy_private_devices %>
ProtectHome=true
ProtectSystem=full
ReadWriteDirectories=<%= @caddy_home %>

[Install]
WantedBy=multi-user.target
